# resources/ledmatrix_widgets.service
[Unit]
Description=ledmatrix_widgets startup script

[Service]
ExecStart=/usr/bin/ledmatrix_widgets

[Install]
WantedBy=multi-user.target